# package/my/luci-app-bbr-mod/root/etc/uci-defaults/bbr
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@bbr[-1]
	add ucitrack bbr
	set ucitrack.@bbr[-1].init=bbr
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
